# buscachat-python/deploy/buscachat-python.service
[Unit]
Description=Buscachat Python API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lozabot
Group=lozabot
WorkingDirectory=/home/lozabot/buscachat-venezuela/buscachat-python
EnvironmentFile=/home/lozabot/buscachat-venezuela/buscachat-python/.env.production
ExecStart=/home/lozabot/.local/bin/uv run python -m fastapi run main.py --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
